# launchd/com.nicklaude.claude-menubar.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nicklaude.claude-menubar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/code/claude-assistant/bin/claude-menubar</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/USERNAME/code/claude-assistant/logs/menubar.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/code/claude-assistant/logs/menubar.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/USERNAME</string>
    </dict>
</dict>
</plist>
